# templates/default/graylog.collector.redhat-init.erb
#! /bin/sh
#
# graylog-collector Starts/stop the "graylog-collector" daemon
#
# chkconfig:   - 95 5
# description: Runs the graylog-collector daemon

### BEGIN INIT INFO
# Provides:          graylog-collector
# Required-Start:    $network $named $remote_fs $syslog
# Required-Stop:     $network $named $remote_fs $syslog
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: Graylog Collector
# Description:       Graylog Collector - Graylog's open source log shipper
### END INIT INFO

# Source function library.
. /etc/rc.d/init.d/functions

RETVAL=0
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin
DESC="Graylog Collector"
NAME=graylog-collector
COLLECTOR_BIN=/usr/share/graylog-collector/bin/graylog-collector
COLLECTOR_ARGS="run -f /etc/graylog/collector/collector.conf"
JAVA=<%= node.graylog2[:collector][:java_bin] %>
JAVA_HOME=<%= node.graylog2[:collector][:java_home] %>
PID_DIR=/var/run/graylog-collector
PID_FILE=$PID_DIR/$NAME.pid
SCRIPTNAME=/etc/init.d/$NAME
LOCKFILE=/var/lock/subsys/$NAME
COLLECTOR_USER=graylog-collector

# Pull in sysconfig settings
[ -f /etc/sysconfig/${NAME} ] && . /etc/sysconfig/${NAME}

# Exit if the package is not installed
[ -e "$COLLECTOR_BIN" ] || exit 0
[ -x "$JAVA" ] || exit 0

getPid() {
    ps auxf | grep -v grep | grep $NAME | awk '{print $2}'
}

start() {
    echo -n $"Starting ${NAME}: "
    install -d -m 755 -o $COLLECTOR_USER -g $COLLECTOR_USER -d $PID_DIR
    daemon --check $JAVA --pidfile=${PID_FILE} --user=${COLLECTOR_USER} \
        "$COLLECTOR_BIN $COLLECTOR_ARGS >> /var/log/graylog-collector/console.log &"
    RETVAL=$?
    echo `getPid` > $PID_FILE
    sleep 2
    [ $RETVAL = 0 ] && touch ${LOCKFILE}
    echo
    return $RETVAL
}

stop() {
    echo -n $"Stopping ${NAME}: "
    killproc -p ${PID_FILE} -d 10 $JAVA
    RETVAL=$?
    [ $RETVAL = 0 ] && rm -f ${PID_FILE} && rm -f ${LOCKFILE}
    echo
    return $RETVAL
}

case "$1" in
    start)
        start
        ;;
    stop)
        stop
        ;;
    status)
        status -p ${PID_FILE} $NAME
        RETVAL=$?
        ;;
    restart|force-reload)
        stop
        start
        ;;
    *)
        N=/etc/init.d/${NAME}
        echo "Usage: $N {start|stop|status|restart|force-reload}" >&2
        RETVAL=2
        ;;
esac

exit $RETVAL
